# installer/linux/DEBIAN/prerm
#!/bin/bash
systemctl stop stt-watchdog || true
systemctl disable stt-watchdog || true
