# stt-watchdog.service
[Unit]
Description=STT Watchdog (manages STT with crash recovery and auto-updates)
After=network.target

[Service]
Type=simple
User=YOUR_USER

# ── Compiled binary install ────────────────────────────────────────────────
# ExecStart=/opt/STT/STT-Watchdog --headless

# ── Python source install ──────────────────────────────────────────────────
# Uncomment and adjust paths for your install:
WorkingDirectory=/home/ai/STT
Environment="PYTHONUNBUFFERED=1"
Environment="HOME=/home/ai"
ExecStart=/home/ai/STT/.venv/bin/python3 /home/ai/STT/watchdog.py --headless
# systemd restarts the watchdog itself if it crashes; the watchdog handles STT restarts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stt-watchdog

[Install]
WantedBy=multi-user.target
