# deploy/com.stt.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stt.watchdog</string>

    <!-- Compiled binary install: replace the two lines below with just the binary path -->
    <!-- <string>/path/to/STT-Watchdog</string>                                        -->
    <key>ProgramArguments</key>
    <array>
        <string>INSTALL_DIR/.venv/bin/python3</string>
        <string>INSTALL_DIR/stt/watchdog.py</string>
        <string>--headless</string>
    </array>

    <key>WorkingDirectory</key>
    <string>INSTALL_DIR</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>INSTALL_DIR/logs/watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>INSTALL_DIR/logs/watchdog.log</string>
</dict>
</plist>
